# systemd/asu-server.service
[Unit]
Description=ASU Server (podman-compose)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=asu
WorkingDirectory=/home/asu/asu-server
# /tmp is tmpfs; recreate the firmware-store bind source on every start
# so podman doesn't mount onto a missing path after a reboot.
ExecStartPre=/usr/bin/mkdir -p /tmp/asu-public-data
ExecStartPre=/usr/bin/logger -t asu-server -p daemon.info "firmware store: host /tmp/asu-public-data -> container /app/public (server+worker)"
ExecStart=/usr/bin/podman-compose up
ExecStop=/usr/bin/podman-compose down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
