# systemd/asu-healthcheck.service
[Unit]
Description=ASU Server Health Check

[Service]
Type=oneshot
User=asu
ExecStart=/home/asu/asu-server/healthcheck.sh
